# gunicorn/tsum20_fastapi.socket
[Unit]
Description=tsum20 socket

[Socket]
ListenStream=/run/tsum20_fastapi.sock

[Install]
WantedBy=sockets.target
